# components/cron/periodic/daily/metarefresh
#!/bin/sh
curl --silent "http://samlbridge:80/samlbridge/module.php/cron/cron.php?key=$SAMLBRIDGE_CRON&tag=daily"
